# redis.ini
[program:redis]
command=/usr/bin/redis-server  /etc/redis.conf
user = root
autostart=true
autorestart=true
startretries=5
stderr_logfile=/log/supervisord/redis_err.log
stdout_logfile=/log/supervisord/redis_std.log