# sshd.ini
[program:sshd]
command=/usr/sbin/sshd -D  -e
user = root
autostart=true
autorestart=true
startretries=5
priority=5
stderr_logfile=/log/supervisord/sshd_err.log
stdout_logfile=/log/supervisord/sshd_out.log